# stage1/units/units/halt.target
[Unit]
Description=Halt
DefaultDependencies=no
Requires=reaper.service
After=reaper.service
AllowIsolate=yes
