# stage1/units/units/poweroff.target
[Unit]
Description=Poweroff
DefaultDependencies=no
Requires=reaper.service
After=reaper.service
AllowIsolate=yes
