# stage1/units/units/reaper.service
[Unit]
Description=rkt apps reaper
AllowIsolate=true
DefaultDependencies=false

[Service]
ExecStart=/reaper.sh
