# stage1/default.target
[Unit]
Description=Rocket stage2 target
DefaultDependencies=false
