# stage1/shutdown.service
[Unit]
Description=Rocket shutdown service

[Service]
ExecStart=/usr/bin/systemctl halt --force
