# stage1/rootfs/aggregate/units/default.target
[Unit]
Description=rkt apps target
DefaultDependencies=false
