# stage1/rootfs/aggregate/units/default.target
[Unit]
Description=Rocket apps target
DefaultDependencies=false
