# stage1/rootfs/aggregate/units/reaper.service
[Unit]
Description=Rocket apps reaper
AllowIsolate=true
DefaultDependencies=false

[Service]
ExecStart=/reaper.sh
