# dist/init/systemd/rkt-gc.service
[Service]
Type=oneshot
ExecStart=/usr/bin/rkt gc --grace-period=5m
