# dist/init/systemd/rkt-gc.timer
[Timer]
OnActiveSec=0s
OnUnitActiveSec=5m

[Install]
WantedBy=default.target
